# etc/ha_healthchecker/ha_healthchecker.service
[Unit]
Description=OpenLab HA Healthchecker

[Service]
ExecStart=/usr/local/bin/ha_healthchecker
