# etc/ha_healthchecker/ha_healthchecker.service
[Unit]
Description=HA Healthchecker

[Service]
ExecStart=/etc/openlab/ha_healthchecker/ha_healthchecker.sh
