# etc/keepalived/healthchecker.service
[Unit]
Description=Healthchecker

[Service]
ExecStart=/etc/healthchecker/healthchecker.sh
