# etc/ha_service/ha_service.timer
[Unit]
Description=OpenLab HA Timer Service

[Timer]
OnActiveSec=10s
OnUnitActiveSec=1m
