# etc/ha_service/ha_service.service
[Unit]
Description=OpenLab HA Service
After=syslog.target network.target

[Service]
ExecStart=/usr/bin/python3 /etc/openlab/kp_local_loop.py
